# system-power/arua_powerhub_indiserver.service
[Unit]
Description=INDI Server Power HUB 
After=multi-user.target

[Service]
Type=idle
Restart=always
RestartSec=5
User=powerhub
#Environment=INDISKEL=/usr/share/indi/arua_powerhub_sk.xml
ExecStart=/usr/local/bin/arua_powerhub_indiserver.sh

[Install]
WantedBy=multi-user.target
